# files/init.d/autoserver
#!/bin/sh /etc/rc.common
# zuoxuzhinan's teaching sucks

START=99

start(){
	/bin/server &	
}

stop(){
	killall server
}

restart(){
	stop
	sleep 1
	start
}
